# profile.d/flutter.sh
export FLUTTER_HOME=/opt/environments/flutter
export DART_SDK=$FLUTTER_HOME/cache/dart-sdk/
export PATH="$PATH:$FLUTTER_HOME/bin:$HOME/.pub-cache/bin"

