# profile.d/poetry_manual.sh
export POETRY_HOME="/opt/environments/poetry"

export PATH="$PATH:$POETRY_HOME/bin"
